# ansible/roles/smarttest/templates/smarttest-@.service
# {{ ansible_managed }}


[Unit]
Description=Launch %j SMART self-test on /dev/%i


[Service]
Type=simple
ExecStart=/usr/sbin/smartctl -t %j /dev/%i
