# ansible/roles/docker_compose/templates/docker-cleanup.timer
# {{ ansible_managed }}

[Unit]
Description=Docker cleanup timer

[Timer]
OnUnitInactiveSec=25h

[Install]
WantedBy=timers.target
