# ansible/roles/systemctl_wait/templates/after.service
# {{ ansible_managed }}
#
# Break circular dependencies with systemctl-wait

[Unit]
After={{ systemctl_wait_after[item]|join(" ") }}
