# ansible/roles/canary/templates/canary.timer
[Unit]
Description=Check for power outages regularly


[Timer]
OnBootSec=5min
OnUnitActiveSec=1min


[Install]
WantedBy=timers.target
